# kukoin-balace.service
[Unit]
Description=test
After=network.target
StartLimitInterval=30
StartLimitBurst=5

[Service]
EnvironmentFile=/root/kukoin_experiments/.env
Type=simple
User=root
ExecStart=/root/.nvm/versions/node/v18.10.0/bin/node /root/kukoin_experiments/balancesTracker.js
Restart=on-failure
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
